# infra/systemd/backend.service
[Unit]
Description=Galactic Conquest Backend (Nodemon)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/bone/oldschoolempire/server
ExecStart=/usr/bin/npm run dev
Restart=always
RestartSec=3
Environment=PATH=/usr/bin:/usr/local/bin
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
