# infra/systemd/prod_backend.service
[Unit]
Description=Dread Horizon Production Backend (Stable)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/bone/oldschoolempire/server
# Ensures we run the built production code
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=3
Environment=PATH=/usr/bin:/usr/local/bin
# Port 3000 is for Production/Beta Testers
Environment=PORT=3000
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
